# aa/aa-quota-remaining-prometheus.service.in
[Unit]
Description=Get bytes remaining for current AA line and expose them to prometheus
After=network.target

[Service]
Type=simple
ExecStart=@PREFIX@/bin/aa-quota-remaining-prometheus

[Install]
WantedBy=network.target
